# projects/monitoring/playbooks/roles/monitoring/template/node_executor.service.j2
[Unit]
Description=Service for node exporter to expose /metrics for prometheus

[Service]
User=monitor
ExecStart=/usr/bin/node_exporter

[Install]
WantedBy=multi-user.target