# projects/monitoring/playbooks/roles/monitoring/templates/node_exporter.service.j2
[Unit]
Description=Service for node exporter to expose /metrics for prometheus

[Service]
User=monitor
ExecStart="{{node_exporter_symlink}}"

[Install]
WantedBy=multi-user.target
